# build/packer/scripts/webserver/cdd-webserver.conf
#!/bin/bash
# upstart job for cdd webserver

description "CDD webserver"
author "CDD"

start on (filesystem and networking) or runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 5 30

script
  # Load environment variable config
  . /etc/default/cdd-webserver

  exec start-stop-daemon -S -c cdd:cdd -x /usr/bin/node /srv/cdd/webserver/server.js -- -L info daemon start
end script
